# deploy/specialwash-backup.service
[Unit]
Description=SpecialWash SQLite daily backup
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/specialwash
Environment=KEEP_DAILY=60
ExecStart=/usr/bin/env bash /root/specialwash/deploy/backup_db.sh

